# EmberOS/src/modules/kiosk/filesystem/pi/.config/autostart/ember-kiosk.desktop
[Desktop Entry]
Name=EmberDefaultKiosk
Type=Application
Exec=/usr/bin/ember-kiosk-launch.sh http://localhost
Terminal=false
