# EmberOS/src/modules/kiosk/filesystem/pi/.config/autostart/nm-applet.desktop
[Desktop Entry]
Name=nm-applet
Type=Application
Exec=nm-applet
Terminal=false
