# EmberOS/src/modules/kiosk/filesystem/pi/.config/autostart/unclutter.desktop
[Desktop Entry]
Name=Unclutter
Type=Application
Exec=unclutter
Terminal=false
