# EmberOS/src/modules/rtc/filesystem/systemd/rtcsync.timer
[Unit]
Description=Run rtcsync.service every 10 minutes to keep the clock synced

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target 