# EmberOS/src/modules/embedpi/filesystem/systemd/bluetooth-autoconnect.service
[Unit]
Description=Bluetooth autoconnect service
Before=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bluetooth-autoconnect.py -d

[Install]
WantedBy=bluetooth.service